# meta-canopy/meta-hpe/meta-gxp/recipes-core/udev/udev-gxp-i2c-passthrough/80-gxp-i2c-passthrough.rules
# Restart HwmonTempSensor when I2C passthrough state changes so that
# newly probed hwmon devices (sbtsi, sbrmi, spd5118) are picked up.
ACTION=="change", SUBSYSTEM=="platform", DRIVER=="gxp-i2c-passthrough", \
	RUN+="/bin/systemctl restart xyz.openbmc_project.hwmontempsensor.service"
